# packaging/tvheadend-manager.service
[Unit]
Description=TvHeadend Manager
Documentation=https://github.com/hiper25/TvHeadend-client-web
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=tvheadend-manager
StateDirectoryMode=0700
WorkingDirectory=/var/lib/tvheadend-manager
EnvironmentFile=-/etc/tvheadend-manager.env
Environment=TVHMON_DATA_DIR=/var/lib/tvheadend-manager
ExecStart=/usr/local/bin/tvheadend-manager
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
UMask=0077

NoNewPrivileges=yes
AmbientCapabilities=
CapabilityBoundingSet=
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RemoveIPC=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
